# root/etc/services.d/plexpy/run
#!/usr/bin/with-contenv bash

cd /app/flextv
exec s6-setuidgid abc python3 ./FlexConnect.py --datadir /config
